# tree/etc/init.d/elive-distro-upgrader
#!/bin/bash
### BEGIN INIT INFO
# Provides:          elive-distro-upgrader
# Required-Start:    $remote_fs $syslog $network
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Perform Debian distribution upgrade on boot
# Description:       This script performs a one-time Debian distribution
#                    upgrade on boot, configured via /etc/default/elive-distro-upgrade.
### END INIT INFO

. /lib/lsb/init-functions

DEFAULTS_FILE="/etc/default/elive-distro-upgrade"
BOOT_UPGRADER="/usr/lib/elive-upgrader/elive-boot-upgrader"
LOCK_FILE="/tmp/.debian-upgrader.lock"

# Read configuration
if [ -r "$DEFAULTS_FILE" ]; then
    . "$DEFAULTS_FILE"
fi

UPGRADE_ENABLED="${UPGRADE_ENABLED:-no}"
UPGRADE_IN_PROGRESS="${UPGRADE_IN_PROGRESS:-no}"
TARGET_CODENAME="${TARGET_CODENAME:-}"
REBOOT_COUNT="${REBOOT_COUNT:-0}"
CURRENT_STEP="${CURRENT_STEP:-}"

case "$1" in
    start)
        log_daemon_msg "Starting Debian distribution upgrade" "elive-distro-upgrader"

        if [ ! -x "$BOOT_UPGRADER" ]; then
            log_failure_msg "$BOOT_UPGRADER not found or not executable"
            exit 1
        fi

        # Run the boot upgrader script
        "$BOOT_UPGRADER"
        exit_code=$?

        if [ "$exit_code" -ne 0 ]; then
            log_end_msg "$exit_code"
        else
            log_end_msg 0
        fi

        exit "$exit_code"
        ;;

    stop)
        log_action_msg "Stopping elive-distro-upgrader..."
        if [ -f "$LOCK_FILE" ]; then
            lock_pid="$(cat "$LOCK_FILE" 2>/dev/null)"
            if [ -n "$lock_pid" ] && kill -0 "$lock_pid" 2>/dev/null; then
                log_action_msg "Sending SIGTERM to upgrade process (PID: $lock_pid)..."
                kill -TERM "$lock_pid" 2>/dev/null || true
                sleep 2
                if kill -0 "$lock_pid" 2>/dev/null; then
                    log_action_msg "Force killing upgrade process (PID: $lock_pid)..."
                    kill -9 "$lock_pid" 2>/dev/null || true
                fi
            fi
            rm -f "$LOCK_FILE"
        fi
        log_action_msg "Upgrade process stopped. Run 'start' to resume or use debian-upgrader --cancel to abort."
        ;;

    restart|force-reload)
        log_action_msg "Restarting elive-distro-upgrader..."
        "$0" stop
        sleep 1
        "$0" start
        ;;

    status)
        log_action_msg "Debian upgrade status:"
        log_action_msg "  Enabled: $UPGRADE_ENABLED"
        log_action_msg "  In Progress: $UPGRADE_IN_PROGRESS"
        log_action_msg "  Target: $TARGET_CODENAME"
        log_action_msg "  Reboot Count: $REBOOT_COUNT"
        log_action_msg "  Current Step: ${CURRENT_STEP:-none}"

        # Get current codename
        current_release_full="$(cat /etc/debian_version 2>/dev/null || echo "unknown")"
        current_version_number="$(echo "$current_release_full" | cut -d'.' -f1)"
        current_cn=""
        case "$current_version_number" in
            15|15.*)  current_cn="duke" ;;
            14|14.*)  current_cn="forky" ;;
            13|13.*)  current_cn="trixie" ;;
            12|12.*)  current_cn="bookworm" ;;
            11|11.*)  current_cn="bullseye" ;;
            10|10.*)  current_cn="buster" ;;
            *)        current_cn="unknown" ;;
        esac

        log_action_msg "  Current Version: $current_cn"
        if [ -f "$LOCK_FILE" ]; then
            lock_pid="$(cat "$LOCK_FILE" 2>/dev/null)"
            if [ -n "$lock_pid" ] && kill -0 "$lock_pid" 2>/dev/null; then
                log_action_msg "  Running: yes (PID: $lock_pid)"
            else
                log_action_msg "  Running: no (stale lock file)"
            fi
        else
            log_action_msg "  Running: no"
        fi
        if [ "$UPGRADE_ENABLED" = "yes" ] && [ "$UPGRADE_IN_PROGRESS" = "yes" ] && [ -n "$TARGET_CODENAME" ]; then
            log_action_msg "Upgrade is configured to run on next boot."
            exit 0
        else
            log_action_msg "Upgrade is not configured to run on next boot."
            exit 3
        fi
        ;;

    *)
        echo "Usage: $0 {start|stop|restart|force-reload|status}" >&2
        exit 3
        ;;
esac

exit 0
